# system/schedule/scheduler.cron
* * * * * /usr/bin/php /var/www/omaya/artisan schedule:run >> /dev/null 2>&1
